# sys-auth/python-validity/files/python3-validity
#!/sbin/openrc-run
command="/usr/lib/python-validity/dbus-service"
command_background=true
pidfile="/var/run/${RC_SVCNAME}.pid"